# jenkinsDocker/config/supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins]
command=/bin/bash -c "exec /etc/init.d/jenkins start"


